# system/on-monitor-change/etc/udev/rules.d/99-on-monitor-change.rules
ACTION=="change", SUBSYSTEM=="drm", RUN+="/usr/bin/systemctl --user start on-monitor-change.service"
