# arch-files/etc/systemd/system/system-backup.service
[Unit]
Description=System backup service

[Service]
Type=simple
ExecStart=/home/bas/.local/bin/system-backup

[Install]
WantedBy=multi-user.target
