# system/system-backup/etc/systemd/system/backup-system.timer
[Unit]
Description=Run backup-system at boot

[Timer]
OnBootSec=1min
Unit=backup-system.service

[Install]
WantedBy=timers.target
