# home/.config/systemd/user/backup-system-notify.path
[Unit]
Description=Watch system backup status

[Path]
PathChanged=/var/lib/backup-system/status

[Install]
WantedBy=default.target
