# home/.config/systemd/user/backup-system-notify.service
[Unit]
Description=Send system backup notification

[Service]
Type=oneshot
ExecStart=/usr/bin/env bash %h/.local/bin/backup-system-notify
