# home/.config/systemd/user/on-monitor-change.service
[Unit]
Description=Monitor change handler

[Service]
Type=oneshot
ExecStart=/usr/bin/fish /home/bas/.local/bin/on-monitor-change.fish
